# systemd/user/end-work-network.service
[Unit]
Description=enable pihole network after work

[Service]
Type=oneshot
ExecStartPre=nmcli conn up wg0
ExecStart=nmcli conn up IndianWay
