# systemd/user/kill-teams.service
[Unit]
Description=Kill Teams after work (5:30 PM)

[Service]
ExecStart=killall teams
